# files/crontab
* 1 * * * /autopgsqlbackup.sh > /backup
# Mandatory blank line